# usr/lib/systemd/system/systemd-resolved.service.d/40_kicksecure-network-conf.conf
## Copyright (C) 2016 - 2025 ENCRYPTED SUPPORT LLC <adrelanos@whonix.org>
## See the file COPYING for copying conditions.

#### meta start
#### project Kicksecure
#### category dns and networking
#### description
## Disable <code>systemd-resolved</code> unless file
## <code>/etc/dns-enable</code> exists.
#### meta end

[Unit]
ConditionPathExists=/etc/dns-enable
