# lib/systemd/system/dnscrypt-proxy.service.d/30_kicksecure.conf
## Copyright (C) 2019 - 2019 ENCRYPTED SUPPORT LP <adrelanos@riseup.net>
## See the file COPYING for copying conditions.

[Service]
ExecStart=
ExecStart=/usr/sbin/dnscrypt-proxy -config /etc/dnscrypt-proxy/kicksecure.toml
